# scripts/profile.d/aliases.sh
alias ll="ls -alF"
alias la="ls -A"
alias l="ls -CF"
bak() { cp "$1"{,.bak};}
mcd() { mkdir -p "$1"; cd "$1";}

